# deploy/supervisor/ai-service.conf
[program:ai-service]
process_name=%(program_name)s
command=/var/www/mock-interview-pro/ai-service/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8001 --workers 2
directory=/var/www/mock-interview-pro/ai-service
autostart=true
autorestart=true
user=www-data
environment=AI_USE_MOCK="false",OLLAMA_URL="http://127.0.0.1:11434",WHISPER_MODEL="base"
redirect_stderr=true
stdout_logfile=/var/log/supervisor/ai-service.log
